# playbooks/roles/lainlet/templates/lainlet.cron.j2
10 5 * * * root /bin/systemctl restart lainlet 2>&1 > /dev/null