# playbooks/roles/bootstrap-etcd/templates/etcd.cron.j2
1 7,12,20,23 * * * root /bin/systemctl restart etcd 2>&1 > /dev/null